# com.google-business-tracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.google-business-tracker</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/helenesteiner/google business tracker/dist/cli.js</string>
        <string>run-daily-check</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>16</integer>
        <key>Minute</key>
        <integer>10</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/helenesteiner/google business tracker/data/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/helenesteiner/google business tracker/data/launchd.error.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/helenesteiner/google business tracker</string>
</dict>
</plist>

